# urirdp-docker/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:urisys-rdp]
command=urisys-rdp --packs rdp,kvm,him,ocr,llm,shell,env,browser --config %(ENV_URISYS_CONFIG)s --events /opt/urirdp/data/events.jsonl serve --host 0.0.0.0 --port %(ENV_URISYS_RDP_PORT)s
directory=/opt/urirdp
user=urisys
environment=URISYS_CONFIG="%(ENV_URISYS_CONFIG)s",URISYS_ALLOW_REAL="%(ENV_URISYS_ALLOW_REAL)s",URISYS_RDP_PORT="%(ENV_URISYS_RDP_PORT)s",URISYS_ENV_POLICY="%(ENV_URISYS_ENV_POLICY)s",OPENROUTER_API_KEY="%(ENV_OPENROUTER_API_KEY)s",LLM_MODEL="%(ENV_LLM_MODEL)s",LLM_BASE_URL="%(ENV_LLM_BASE_URL)s",LLM_TEMPERATURE="%(ENV_LLM_TEMPERATURE)s",LLM_MAX_TOKENS="%(ENV_LLM_MAX_TOKENS)s"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
